# etc/periodic/daily/Makefile
# $FreeBSD$

FILES=	daily \
	    daily.conf

.include <bsd.prog.mk>
