# etc/periodic/monthly/Makefile
# $FreeBSD$

FILES=	monthly \
	    monthly.conf

.include <bsd.prog.mk>
