# etc/periodic/security/Makefile
# $FreeBSD$

FILES=	security\
        security.conf

.include <bsd.prog.mk>
