# etc/periodic/weekly/Makefile
# $FreeBSD$

FILES=	weekly\
        weekly.conf

.include <bsd.prog.mk>
